# aaa.conf
[program:aaa]
command=python3 /home/curator2/aaa/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/aaa/test.err.log
stdout_logfile=/home/curator2/aaa/test.out.log